# nesop-store-deployment-20250715_154232/nesop-store.service
[Unit]
Description=NESOP Store Application
After=network.target
Wants=network.target

[Service]
Type=notify
User=nesop
Group=nesop
WorkingDirectory=/Users/kasonsuchow/Documents/github/NESOP-Store/nesop-store-deployment-20250715_154232
Environment=PATH=/Users/kasonsuchow/Documents/github/NESOP-Store/nesop-store-deployment-20250715_154232/venv/bin
ExecStart=/Users/kasonsuchow/Documents/github/NESOP-Store/nesop-store-deployment-20250715_154232/venv/bin/gunicorn --bind 0.0.0.0:8080 --workers 4 --timeout 120 --access-logfile logs/access.log --error-logfile logs/error.log wsgi:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
RestartSec=10
Restart=always

[Install]
WantedBy=multi-user.target
